# 27-backup/files/borg-backup.timer
[Unit]
Description=Borg Backup

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
